# .docker/php-fpm/etc/supervisor.conf
[supervisord]
nodaemon=true

[program:php7.0-fpm]
command=/etc/init.d/php7.0-fpm start